# backend/sarjet-backend.service
[Unit]
Description=SarjEt Backend API
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/your/SarjEt/backend
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sarjet-backend

[Install]
WantedBy=multi-user.target
